# roles/machines/server/service-auto-system-update/files/auto-system-update.timer
[Unit]
Description=Weekly automatic system update

[Timer]
OnCalendar=Sat *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
