# roles/machines/server/files/rclone-sync.timer
[Unit]
Description=Weekly rclone backup sync

[Timer]
OnCalendar=Sun *-*-* 05:00:00

[Install]
WantedBy=rclone-sync.target
