# roles/machines/server/service-venmo-lunchmoney-ai/files/venmo-lunchmoney-ai.timer
[Unit]
Description=Automatic venmo Reconciliation in lunchmoney

[Timer]
OnCalendar=*-*-* *:00,30:00

[Install]
WantedBy=timers.target
