# roles/machines/server/files/systemd/auto-system-update.timer
[Unit]
Description=Weekly automatic system update

[Timer]
OnCalendar=Sun *-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
