# roles/machines/server/files/systemd/rclone-cryptcheck.timer
[Unit]
Description=Monthly backup encryption verification
After=rclone-sync.service

[Timer]
OnCalendar=*-*-01 05:00:00

[Install]
WantedBy=timers.target
