# roles/machines/server/files/systemd/record-file-history.timer
[Unit]
Description=Weekly file history record

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target
