# roles/machines/server/files/systemd/speedtest-monitor.timer
[Unit]
Description=Twice daily speedtest

[Timer]
OnCalendar=*-*-* 8,22:00:00

[Install]
WantedBy=timers.target
