# roles/machines/server/files/systemd/venmo-auto-cashout.timer
[Unit]
Description=Automatic venmo cashout every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
