# roles/pacsync/files/pacsync.timer
[Unit]
Description=Run pacman database sync hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
