# roles/pacsync/files/pacsync.timer
[Unit]
Description=Hourly pacman database sync

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
