# roles/machines/server/service-codex/files/codex.service
[Unit]
Description=Codex app server
After=network-online.target import-environment.service
Wants=network-online.target
Requires=import-environment.service

[Service]
Type=simple
Environment=SSH_AUTH_SOCK=/run/ssh-agent-proxy.sock
Environment=CODEX_HOME=/home/evan/.config/codex
EnvironmentFile=-/home/evan/.config/codex/mcp.env
ExecStart=/usr/local/libexec/codex/app-server
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
