# roles/machines/server/service-bambulab-lights-off/files/bambulab-lights-off.timer
[Unit]
Description=Turn off Bambu Lab printer lights when idle

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
