# roles/machines/server/service-iptorrents-keepalive/files/iptorrents-keepalive.service
[Unit]
Description=Run a periodic IPTorrents authenticated request

[Service]
Type=oneshot
User=evan
Environment=HOME=/home/evan
ExecStart=/usr/bin/ipt search "linux" --limit 1
