# roles/machines/server/service-opencode/files/opencode.service
[Unit]
Description=OpenCode server
After=network-online.target import-environment.service
Wants=network-online.target
Requires=import-environment.service

[Service]
Type=simple
Environment=SSH_AUTH_SOCK=/run/ssh-agent-proxy.sock
Environment=OPENCODE_HOST=http://localhost:9995
EnvironmentFile=-/home/evan/.config/opencode/mcp.env
ExecStart=/usr/local/libexec/opencode/serve
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
