# roles/machines/server/service-opencode/files/opencode-ssh-agent-proxy@.service
[Unit]
Description=SSH agent proxy connection (1Password via MacBook)

[Service]
Type=simple
ExecStart=/usr/local/bin/opencode-ssh-agent-proxy
StandardInput=socket
StandardOutput=socket
StandardError=journal
# Each instance handles one agent request then exits
TimeoutStopSec=20
